# hashicorp-vault/ssh-host-certificate-renew.timer
[Unit]
Description=SSH Host Certificate Renewer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
